# userspace/inittab
# [DistroName]OS inittab
# BusyBox init configuration

# System initialization
::sysinit:/etc/init.d/rcS

# Start getty on console
::respawn:/sbin/getty -L console 115200 vt100
tty1::respawn:/sbin/getty 38400 tty1
tty2::respawn:/sbin/getty 38400 tty2

# Stuff to do before rebooting
::shutdown:/bin/echo "Shutting down..."
::shutdown:/bin/umount -a -r
::shutdown:/sbin/swapoff -a

# Ctrl-Alt-Del
::ctrlaltdel:/sbin/reboot
